# deploy/synaplan-watchguard.timer
[Unit]
Description=Run Synaplan watchguard every 2 minutes
Documentation=https://github.com/metadist/synaplan-website

[Timer]
# 60s after boot to give docker/network a chance to settle, then poll every
# 2 minutes. Adjust OnUnitActiveSec if GHCR rate-limits or you want to slow
# the loop down.
OnBootSec=60
OnUnitActiveSec=2min
Unit=synaplan-watchguard.service
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
