# deploy/synaplan-watchguard.service
[Unit]
Description=Synaplan website watchguard (pulls ghcr image, re-up compose if changed)
Documentation=https://github.com/metadist/synaplan-website
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/synaplan-watchguard.sh
# Bound long-running deploys so the timer doesn't accumulate failed runs.
TimeoutStartSec=300
Nice=10

[Install]
WantedBy=multi-user.target
